# deploy/polysignal.service
[Unit]
Description=PolySignal Trading Bot
After=network.target

[Service]
Type=simple
User=polysignal
Group=polysignal
WorkingDirectory=/home/ninetheme-bot/htdocs/bot.ninetheme.com/polysignal/polysignal
ExecStart=/usr/bin/python3 main.py
ProtectSystem=strict
ReadWritePaths=/home/ninetheme-bot/htdocs/bot.ninetheme.com/polysignal/polysignal/data
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
